# wyoming-openclaw.service
# Wyoming-OpenClaw systemd service
#
# Setup:
#   1. Copy to /etc/systemd/system/wyoming-openclaw.service
#   2. Edit User, WorkingDirectory, and ExecStart paths
#   3. Create /etc/wyoming-openclaw.env with your credentials (see below)
#   4. systemctl daemon-reload && systemctl enable --now wyoming-openclaw
#
# /etc/wyoming-openclaw.env should contain:
#   GATEWAY_TOKEN=your_gateway_token
#   HA_TOKEN=your_ha_token          # optional
#   HA_URL=http://homeassistant:8123 # optional

[Unit]
Description=Wyoming OpenClaw Bridge
After=network.target

[Service]
Type=simple
User=wyoming
WorkingDirectory=/opt/wyoming-openclaw
EnvironmentFile=/etc/wyoming-openclaw.env
ExecStart=/opt/wyoming-openclaw/venv/bin/python wyoming_openclaw.py \
    --port 10600 \
    --gateway-url http://127.0.0.1:18789 \
    --session-id voice-assistant
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
