# wyoming-openclaw.service
[Unit]
Description=Wyoming OpenClaw Bridge
After=network.target

[Service]
Type=simple
User=your-user
WorkingDirectory=/path/to/wyoming-openclaw
Environment="PATH=/usr/local/bin:/usr/bin:/bin"
ExecStart=/path/to/wyoming-openclaw/venv/bin/python wyoming_openclaw.py --port 10600 --gateway-url http://127.0.0.1:18789 --token YOUR_GATEWAY_TOKEN --session-id voice-assistant
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
